# imagebuilder/target/linux/gem6xxx/base-files/etc/hotplug.d/net/10-netctl
#!/bin/sh

if [ "$ACTION" = add ]; then
	[ -f /sbin/smp.sh ] && /sbin/smp.sh
fi
