# imagebuilder/target/linux/gem6xxx/base-files/etc/uci-defaults/90_advanced_system_log_setting
#!/bin/sh

str=$(uci -q get system.@system[0].enable_sensitive_log)

if [[ -z $str ]]; then
    uci set system.@system[0].ap_start='0'
    uci set system.@system[0].gps_start='0'
    uci set system.@system[0].net_start='0'
    uci set system.@system[0].md_start='0'
    uci set system.@system[0].ap_save_mode='1'
    uci set system.@system[0].ap_max_files='1000'
    uci set system.@system[0].ap_max_bootlog_folders='5'
    uci set system.@system[0].ap_auto_logging='0'
    uci set system.@system[0].md_save_mode='1'
    uci set system.@system[0].md_gps_location_log_enable='0'
    uci set system.@system[0].md_level='-1'
    uci set system.@system[0].md_dumpback_mode='0'
    uci set system.@system[0].md_log_size='20'
    uci set system.@system[0].md_cyclic_size='200'
    uci set system.@system[0].md_flush_pls_log='0'
    uci set system.@system[0].log_store_path='/data/debuglog'
    uci set system.@system[0].enable_sensitive_log='0'
    uci set system.@system[0].log_size='2048'
    uci set system.@system[0].md_memdump_wait_time='15'
    uci set system.@system[0].md_log_ip='0.0.0.0'
    uci set system.@system[0].md_log_port='10001'
    uci commit system
fi

exit 0


