# XS.EDONGLEV3.01.01.10-ImageBuilder-20221124/target/linux/gem6xxx/xs5g01_edonglev3/base-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo 3 > /sys/devices/platform/11201000.usb/mode
. /usr/sbin/rtl_phy_patch.sh -p 1
. /usr/sbin/rtl_phy_patch.sh -p 2
. /usr/sbin/rtl_phy_patch.sh -p 5
. /usr/sbin/xs_mac.sh

mkdir -p xnvfile

if [ ! -d "/data/ssl" ] 
then
    echo "Error: Directory /data/ssl does not exists."
    mkdir -p /data/ssl
    /etc/init.d/uhttpd restart
fi

exit 0
